# etc/profile.d/03-register-binfmt.sh
sudo bin/qemu-binfmt-conf.sh >& /dev/null
